# config/systemd/user/mbsync.timer
[Unit]
Description=Synchronize service for mailboxes

[Timer]
OnBootSec=1m
OnUnitActiveSec=4m
Unit=mbsync.service

[Install]
WantedBy=timers.target
