# config/systemd/user/random-wallpaper.service
[Unit]
Description=Random wallpaper

[Service]
Type=oneshot
ExecStart=random-wallpaper
