# config/systemd/user/random-wallpaper.timer
[Unit]
Description=Random wallpaper timer

[Timer]
OnBootSec=30min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target