# config/systemd/user/random-wall.service
[Unit]
Description=Random wallpaper

[Service]
Type=oneshot
ExecStart=/home/sqve/scripts/random-wall.sh
